# deploy/awaria-backup.service
[Unit]
Description=Awaria DB backup to the gcode pendrive

[Service]
Type=oneshot
ExecStart=/usr/local/bin/awaria-backup
User=admin
